# search/supervisord.conf
[supervisord]
nodaemon=true
logfile = /var/log/supervisord/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
 
[program:cron]
autorestart=false
command=cron -f
 
[program:sphinx]
autorestart=false
command=searchd -c /root/sphinx.conf --nodetach
stderr_logfile = /var/log/supervisord/sphinx.log
stdout_logfile = /var/log/supervisord/sphinx.log