# extra/monaba
#!/sbin/runscript
 
path="/home/user/Monaba"
command="/home/user/Monaba/dist/build/Monaba/Monaba"
command_args="config/settings.yml"
user="user"
group="user"
name="monaba"
pidfile="/var/run/monaba.pid"
logfile="/var/log/monaba.log"

depend() {
  need net
  need postgresql-9.4
}

start() {
  ebegin "Starting Monaba"
  if [ ! -f $logfile ]; then 
      touch $logfile
  fi
  chown "$user:$group" $logfile
  start-stop-daemon --start --user $user --group $group --chdir $path --name $name \
    --pidfile $pidfile --make-pidfile --background \
    --exec "/bin/bash" -- -c "exec $command $command_args > $logfile 2>&1"
  eend $?
}

stop() {
  ebegin "Stopping Monaba"
  start-stop-daemon --stop --exec $command --pidfile $pidfile
  eend $?
}
